# systemd/pqc-monitor-web.service
# PQC-Monitor: Web Service Unit
# Runs the Flask application under Gunicorn for production use.
#
# Installation:
#   sudo cp systemd/pqc-monitor-web.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pqc-monitor-web
#
# Logs:
#   journalctl -u pqc-monitor-web -f
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=PQC-Monitor Web Service
Documentation=https://github.com/your-org/pqc-monitor
After=network-online.target
Wants=network-online.target
PartOf=pqc-monitor.target
StartLimitIntervalSec=90
StartLimitBurst=5

[Service]
Type=simple

# ── Identity ─────────────────────────────────────────────────────────────────
User=pqcmonitor
Group=pqcmonitor
WorkingDirectory=/opt/pqc-monitor

# ── Environment ──────────────────────────────────────────────────────────────
# Defaults below are overridden by values in pqc-monitor.env if set there.
# systemd does not support bash ${VAR:-default} syntax in ExecStart;
# declare the defaults here instead.
EnvironmentFile=/etc/pqc-monitor/pqc-monitor.env
Environment=PYTHONUNBUFFERED=1
Environment=PQC_BIND=127.0.0.1:5000
Environment=PQC_WEB_WORKERS=2

# ── Startup ───────────────────────────────────────────────────────────────────
# Gunicorn listens on 127.0.0.1:5000 by default.
# Put nginx or caddy in front for TLS termination.
ExecStart=/opt/pqc-monitor/.venv/bin/gunicorn \
    --bind ${PQC_BIND} \
    --workers ${PQC_WEB_WORKERS} \
    --worker-class sync \
    --worker-tmp-dir /tmp \
    --timeout 120 \
    --keep-alive 5 \
    --access-logfile - \
    --error-logfile - \
    --capture-output \
    --forwarded-allow-ips "*" \
    "app_factory:create_app()"

ExecReload=/bin/kill -HUP $MAINPID

# ── Hardening ─────────────────────────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/pqc-monitor /var/log/pqc-monitor
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=10

# ── Logging ───────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pqc-monitor-web

[Install]
WantedBy=pqc-monitor.target
